# p1nrgy/rootfs/etc/services.d/p1nrgy/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

function _hass {
  ADDON_VERSON=$(jq --raw-output '.version' /app/ha_addon_version.json)
  _info "Home Assistant Add-on release: ${ADDON_VERSON}"
  CHECK_UPDATE=$(curl -s "https://api-check.duckdns.org/dsmr-datalogger-addon/${ADDON_VERSON}")
  if [[ "$CHECK_UPDATE" == *"response_string"* ]]; then
    OUTPUT=$(echo $CHECK_UPDATE | jq --raw-output .response_string)
    _info "$OUTPUT"
  else
    _info "Home Assistant Add-on: Update check failed"
  fi

  declare dsmr_enabled
  declare dsmr_mode
  declare dsmr_serial_device
  declare dsmr_serial_baudrate
  declare dsmr_serial_databits
  declare dsmr_serial_parity
  declare dsmr_serial_stopbits
  declare dsmr_serial_timeout

  CONFIG_PATH=/data/options.yaml

  dsmr_enabled=$(bashio::config 'dsmr_enabled')
  if bashio::var.true "${dsmr_enabled}"; then
    bashio::log.info "DSMR reading is enabled"
  else
    bashio::log.info "DSMR reading is disabled"
  fi

  dsmr_mode=$(bashio::config 'dsmr_mode')
  bashio::log.info "DSMR mode set to: ${dsmr_mode}"
  dsmr_serial_device=$(bashio::config 'dsmr_serial_device')
  bashio::log.info "DSMR serial device set to: ${dsmr_serial_device}"

  dsmr_serial_baudrate=$(bashio::config 'dsmr_serial_baudrate')
  bashio::log.info "DSMR serial baudrate set to: ${dsmr_serial_baudrate}"

  dsmr_serial_databits=$(bashio::config 'dsmr_serial_databits')
  bashio::log.info "DSMR serial databits set to: ${dsmr_serial_databits}"

  dsmr_serial_parity=$(bashio::config 'dsmr_serial_parity')
  bashio::log.info "DSMR serial parity set to: ${dsmr_serial_parity}"

  dsmr_serial_stopbits=$(bashio::config 'dsmr_serial_stopbits')
  bashio::log.info "DSMR serial stopbits set to: ${dsmr_serial_stopbits}"

  dsmr_serial_timeout=$(bashio::config 'dsmr_serial_timeout')
  bashio::log.info "DSMR serial timeout set to: ${dsmr_serial_timeout}"

  export DSMR_ENABLED=dsmr_enabled
  export DSMR_MODE=dsmr_mode
  export DSMR_SERIAL_DEVICE=dsmr_serial_device
  export DSMR_SERIAL_BAUDRATE=dsmr_serial_baudrate
  export DSMR_SERIAL_DATABITS=dsmr_serial_databits
  export DSMR_SERIAL_PARITY=dsmr_serial_parity
  export DSMR_SERIAL_STOPBITS=dsmr_serial_stopbits
  export DSMR_SERIAL_TIMEOUT=dsmr_serial_timeout

  bashio::log.info "p1nrgy started"
}

_hass

## Run your program
exec /usr/bin/p1nrgy
