# p1nrgy/rootfs/etc/services.d/p1nrgy/run
#!/usr/bin/with-contenv bashio
set -euo pipefail

main() {
  local addon_version response message env_file

  addon_version=$(jq --raw-output '.version' /app/ha_addon_version.json)
  bashio::log.info "Home Assistant add-on release: ${addon_version}"

  if response=$(curl -sf "https://api-check.duckdns.org/dsmr-datalogger-addon/${addon_version}"); then
    message=$(jq --raw-output '.response_string // empty' <<< "${response}")
    if [[ -n "${message}" ]]; then
      bashio::log.info "${message}"
    else
      bashio::log.warning "Update check response missing message"
    fi
  else
    bashio::log.warning "Home Assistant add-on update check failed"
  fi

  DSMR_ENABLED="$(bashio::config 'dsmr_enabled')"
  DSMR_MODE="$(bashio::config 'dsmr_mode')"
  DSMR_SERIAL_DEVICE="$(bashio::config 'dsmr_serial_device')"
  DSMR_SERIAL_BAUDRATE="$(bashio::config 'dsmr_serial_baudrate')"
  DSMR_SERIAL_DATABITS="$(bashio::config 'dsmr_serial_databits')"
  DSMR_SERIAL_PARITY="$(bashio::config 'dsmr_serial_parity')"
  DSMR_SERIAL_STOPBITS="$(bashio::config 'dsmr_serial_stopbits')"
  DSMR_SERIAL_TIMEOUT="$(bashio::config 'dsmr_serial_timeout')"

  MQTT_HOST="$(bashio::services mqtt 'host')"
  MQTT_PORT="$(bashio::services mqtt 'port')"
  MQTT_USERNAME="$(bashio::services mqtt 'username')"
  MQTT_PASSWORD="$(bashio::services mqtt 'password')"
  MQTT_CLIENT_ID="p1nrgy_hass_addon"
  MQTT_TOPIC="$(bashio::config 'mqtt_topic')"

  export DSMR_ENABLED
  export DSMR_MODE
  export DSMR_SERIAL_DEVICE
  export DSMR_SERIAL_BAUDRATE
  export DSMR_SERIAL_DATABITS
  export DSMR_SERIAL_PARITY
  export DSMR_SERIAL_STOPBITS
  export DSMR_SERIAL_TIMEOUT

  export MQTT_BROKER="tcp://${MQTT_HOST}:${MQTT_PORT}"
  export MQTT_USERNAME
  export MQTT_PASSWORD
  export MQTT_CLIENT_ID
  export MQTT_TOPIC

  # Server settings
  # Server settings:
  export SERVER_HOST="0.0.0.0"
  export SERVER_PORT=3000
  export SERVER_READ_TIMEOUT=60

  bashio::log.info "p1nrgy started"
}

main

exec /usr/bin/p1nrgy
